# conf.d/dosbox.conf
[program:dosbox]
command=dosbox -conf /app/dosbox/autoexec.conf
autorestart=true

