# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -rfbauth /app/vnc_pass
autorestart=true
